# com.claude-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-relay</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/gaylonvorwaller/claude-relay/server.js</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/Users/gaylonvorwaller/claude-relay/logs/relay-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/gaylonvorwaller/claude-relay</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>RELAY_MESSAGE_RETENTION_DAYS</key>
        <string>7</string>
        <key>RELAY_MESSAGE_MAX_DISK_MB</key>
        <string>100</string>
        <key>RELAY_CACHE_MAX_MESSAGES</key>
        <string>500</string>
        <key>RELAY_CACHE_MAX_MB</key>
        <string>10</string>
        <key>RELAY_LOG_RETENTION_DAYS</key>
        <string>7</string>
        <key>RELAY_LOG_MAX_TOTAL_MB</key>
        <string>50</string>
        <key>RELAY_LOG_MAX_FILE_MB</key>
        <string>10</string>
        <key>RELAY_ADMIN_CLIENT_IDS</key>
        <string>CODEX2,CC2</string>
    </dict>
</dict>
</plist>
